# ssh/hotcoder_tunnel/supervisord.conf
[supervisord]
nodaemon=true

[program:ssh_db]
command=ssh -F /ssh/config -o "StrictHostKeyChecking no" -L 0.0.0.0:3306:localhost:3306 hotcoder-demo-db1